# roles/nft/files/etc/systemd/system/nftables.service.d/override.conf
[Unit]
PropagatesReloadTo=fail2ban.service
