# roles/openvpn/files/etc/systemd/system/openvpnmapping.service
[Unit]
Description=update nft for vpn clients

[Service]
Type=oneshot
ExecStart=/usr/local/bin/gen_nft_defines
