# roles/openvpn/files/etc/systemd/system/openvpnmapping.path
[Unit]
Description=Watchout for new openvpn clients

[Path]
PathChanged=/var/log/openvpn/ipp.txt

[Install]
WantedBy=multi-user.target

