# roles/fail2ban/files/etc/systemd/system/fail2ban.service.d/override.conf
[Unit]
PartOf=nftables.service

[Service]
ExecReload=
ExecReload=/usr/bin/fail2ban-client restart
