# roles/afraid/files/etc/systemd/system/update-afraid.service
[Unit]
Description=update dyn dns record on afraid.org

[Service]
Type=oneshot
ExecStart=/usr/local/bin/update_afraid.sh
